# daemon/pixypy.service
[Unit]
Description=Start pixy python script
After=network.target
StartLimitIntervalSec=0
Requires=pixyc.service

[Service]
Type=simple
Restart=always
RestartSec=1
User=root
ExecStart=/usr/bin/python /home/pi/src/Pixy-SIW/main.py
ExecStartPre=/bin/sleep 3

[Install]
WantedBy=multi-user.target
